# packaging/linux/conf/antz.service
[Unit]
Description=ANT+ Device Discovery Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/antz -d 0 -e 1
Restart=on-failure
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
